# bess-apg-scheduler.service
[Unit]
Description=BESS APG Scheduler - Automatischer Import österreichischer Spot-Preise
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/bess-simulation
ExecStart=/opt/bess-simulation/venv/bin/python /opt/bess-simulation/apg_scheduler_linux.py
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/opt/bess-simulation
Environment=PYTHONUNBUFFERED=1

# Ressourcen-Limits
MemoryMax=512M
CPUQuota=50%

# Restart-Policy
Restart=on-failure
RestartSec=300

# Sicherheit
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bess-simulation/instance /var/log/bess

[Install]
WantedBy=multi-user.target
